# mlnx_ofed/mlnx-ofed-kernel-4.9/ofed_scripts/mlnx_interface_mgr@.service
[Unit]
Description=mlnx_interface_mgr - configure %I

[Service]
Type=simple
ExecStart=/bin/mlnx_interface_mgr.sh %I
ExecStop=/bin/true
